# .cache/crontab/crontab.bak
@reboot cd $HOME/bot && $HOME/.local/bin/honcho start >> $HOME/logs/honcho.log 2>&1
*/10 * * * * pgrep -f honcho > /dev/null || (cd $HOME/bot && $HOME/.local/bin/honcho start >> $HOME/logs/honcho_watchdog.log 2>&1 &)
